# systemd/system/snap-snapd-13640.mount
[Unit]
Description=Mount unit for snapd, revision 13640
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/snapd_13640.snap
Where=/snap/snapd/13640
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
