# systemd/system/snap-amazon\x2dssm\x2dagent-4046.mount
[Unit]
Description=Mount unit for amazon-ssm-agent, revision 4046
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/amazon-ssm-agent_4046.snap
Where=/snap/amazon-ssm-agent/4046
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
